# s6-rc.d/areas_updater/run
#!/usr/bin/env bash

sleep 1
/overpass/bin/rules_loop.sh "/data/db/"
